# projects/S905/packages/aml_linux_fd628/system.d/fd628.service
[Unit]
Description=Amlogic FD628 service
ConditionPathExists=/dev/fd628_dev

[Service]
Type=oneshot
ExecStart=/usr/sbin/fd628
RemainAfterExit=yes

[Install]
WantedBy=basic.target
